# dotfiles/example.service
[Unit]

[Service]
ExecStart=/bin/service

[Install]
WantedBy=multi-user.target
