# dotfiles/example.service
[Service]
Type=oneshot
ExecStart=/bin/service
RemainAfterExit=yes
StandardOutput=tty

[Install]
WantedBy=multi-user.target
