# dotfiles/rc.local
# Place here your custom commands to be run at startup
loadkeys /etc/extrakeys
