# arch/root/rc.local
#!/bin/bash
#Place here the commands that will run at startup
loadkeys /etc/extrakeys
